# openHAB/openhab-mqttclient/etc/systemd/system/mqtt.leds.service
[Service]
Type=simple
ExecStart=/usr/bin/python /media/mqtt.bedroom.led.py &

[Install]
WantedBy=multi-user.target
